# scripts/setup/com.fas.gemini-b.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- FAS Gemini CLI Account B (Cross-verification) launchd plist
     Auto-starts Gemini CLI session B in tmux on login.

     Install:
       cp scripts/setup/com.fas.gemini-b.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.fas.gemini-b.plist

     Uninstall:
       launchctl unload ~/Library/LaunchAgents/com.fas.gemini-b.plist
       rm ~/Library/LaunchAgents/com.fas.gemini-b.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fas.gemini-b</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>tmux new-session -d -s fas-gemini-b 'GEMINI_ACCOUNT=B bash /Users/user/FAS-operations/scripts/gemini/gemini_wrapper.sh b' 2>/dev/null || true</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/user/FAS-operations/logs/gemini-b-launch.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/user/FAS-operations/logs/gemini-b-launch.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/user/FAS-operations</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/user</string>
    </dict>
</dict>
</plist>
