# scripts/setup/com.fas.hunter.plist
<!-- FAS Hunter launchd plist
     Auto-starts hunter tmux session on login.

     Install:
       cp scripts/setup/com.fas.hunter.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.fas.hunter.plist

     Uninstall:
       launchctl unload ~/Library/LaunchAgents/com.fas.hunter.plist
       rm ~/Library/LaunchAgents/com.fas.hunter.plist
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fas.hunter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>/Users/user/fas-hunter-agent/scripts/hunter_watchdog.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/user/Library/Logs/fas-hunter.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/user/Library/Logs/fas-hunter.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/user/.nvm/versions/node/v22.22.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/user</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
